# mira-daemon.service
[Unit]
Description=Mira Code Intelligence Daemon
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=peter
WorkingDirectory=/home/peter/Mira
ExecStart=/home/peter/Mira/scripts/start-daemon.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
