# mira-http.service
# Mira HTTP Server - Remote MCP Access
# Run Mira as an HTTP server for multi-device/multi-session access
#
# Installation:
#   1. Edit this file: replace YOUR_USERNAME with your username
#   2. sudo cp mira-http.service /etc/systemd/system/
#   3. sudo systemctl daemon-reload
#   4. sudo systemctl enable mira-http
#   5. sudo systemctl start mira-http
#
# Configure Claude Code:
#   claude mcp add mira -t http http://localhost:3000/mcp
#
# To set API keys:
#   sudo systemctl edit mira-http
#   Add:
#     [Service]
#     Environment=GEMINI_API_KEY=your-key
#     Environment=MIRA_AUTH_TOKEN=your-token

[Unit]
Description=Mira HTTP Server (MCP over HTTP for remote access)
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/.mira
Environment=DATABASE_URL=sqlite:///home/YOUR_USERNAME/.mira/data/mira.db
Environment=QDRANT_URL=http://localhost:6334
ExecStart=/home/YOUR_USERNAME/.mira/target/release/mira serve-http --port 3000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
