# mira.service
# Mira Service - MCP + Chat + Indexer
# Full Mira server on a single port
#
# Installation:
#   1. Edit this file: replace YOUR_USERNAME with your username
#   2. sudo cp mira.service /etc/systemd/system/
#   3. sudo systemctl daemon-reload
#   4. sudo systemctl enable mira
#   5. sudo systemctl start mira
#
# To set API keys:
#   sudo systemctl edit mira
#   Add:
#     [Service]
#     Environment=DEEPSEEK_API_KEY=your-key
#     Environment=GEMINI_API_KEY=your-key
#     Environment=MIRA_AUTH_TOKEN=your-token
#
# Endpoints:
#   MCP:    http://localhost:3000/mcp
#   Chat:   http://localhost:3000/api/chat/stream
#   Health: http://localhost:3000/health

[Unit]
Description=Mira Server (MCP + Chat + Indexer)
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/Mira
Environment=DATABASE_URL=sqlite:///home/YOUR_USERNAME/.mira/mira.db
Environment=QDRANT_URL=http://localhost:6334
EnvironmentFile=-/home/YOUR_USERNAME/.mira/.env
ExecStart=/home/YOUR_USERNAME/Mira/target/release/mira serve-http --port 3000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
