# mira-server.service
# Mira MCP Server - HTTP Mode
#
# Installation:
#   1. Edit this file and replace YOUR_USERNAME and paths as needed
#   2. sudo cp mira-server.service /etc/systemd/system/
#   3. sudo systemctl daemon-reload
#   4. sudo systemctl enable mira-server
#   5. sudo systemctl start mira-server
#
# To set API keys securely:
#   sudo systemctl edit mira-server
#   Add:
#     [Service]
#     Environment=GEMINI_API_KEY=your-key
#     Environment=MIRA_AUTH_TOKEN=your-token

[Unit]
Description=Mira MCP Server (HTTP)
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/.mira
Environment=DATABASE_URL=sqlite:///home/YOUR_USERNAME/.mira/data/mira.db
Environment=QDRANT_URL=http://localhost:6334
ExecStart=/home/YOUR_USERNAME/.mira/target/release/mira serve-http --port 3000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
