# service/mira.service
[Unit]
Description=Mira Intelligence Layer for Claude Code
Documentation=https://github.com/ConaryLabs/Mira
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/mira
Restart=always
RestartSec=5

# Environment
Environment=MIRA_PORT=3000
# MIRA_LISTEN defaults to 127.0.0.1 (localhost only - secure)
# To expose externally: MIRA_LISTEN=0.0.0.0 (requires auth for /api/*)
# Environment=MIRA_LISTEN=127.0.0.1
# For external access, set allowed CORS origins:
# Environment=MIRA_CORS_ORIGINS=https://claude.ai,https://your-domain.com
EnvironmentFile=%h/.mira/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mira

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.mira %h/Mira/data
PrivateTmp=true

[Install]
WantedBy=default.target
