# service/mira.service
[Unit]
Description=Mira Intelligence Layer for Claude Code
Documentation=https://github.com/ConaryLabs/Mira
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/mira
Restart=always
RestartSec=5

# Environment
Environment=MIRA_PORT=3000
EnvironmentFile=%h/.mira/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mira

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.mira %h/Mira/data
PrivateTmp=true

[Install]
WantedBy=default.target
