# scripts/mira-qdrant.service
[Unit]
Description=Qdrant Vector Database for Mira
After=network.target

[Service]
Type=simple
User=MIRA_USER
Group=MIRA_USER
WorkingDirectory=MIRA_PATH/backend
ExecStart=MIRA_PATH/backend/bin/qdrant --config-path MIRA_PATH/backend/config/config.yaml
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mira-qdrant

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
